# .crontab
0 * * * * /Users/kevin/Developer/bash/growl-time.sh 
